# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xvfb]
command=Xvfb :0 -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx24
autorestart=true

[program:fluxbox]
command=fluxbox
autorestart=true

[program:x11vnc]
command=x11vnc -display :0 -forever -shared -nopw -rfbport 5900
autorestart=true

[program:novnc]
command=/usr/bin/websockify --web /usr/share/novnc 6080 localhost:5900
autorestart=true

[program:opencode-server]
command=opencode serve --port 4096 --hostname 0.0.0.0
autorestart=true
stopwaitsecs=30
startsecs=10
priority=10
stdout_logfile=/app/opencode/logs/opencode-server.log
stderr_logfile=/app/opencode/logs/opencode-server.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
environment=OPENCODE_SERVER_PASSWORD="opencode-dev-2026",OPENCODE_SERVER_USERNAME="opencode"
stopasgroup=true
killasgroup=true

[program:app]
command=/app/opencode/app/start_app.sh
autorestart=true
stopwaitsecs=30
startsecs=15
priority=20
stdout_logfile=/app/opencode/logs/app.log
stderr_logfile=/app/opencode/logs/app.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
stopasgroup=true
killasgroup=true
after=opencode-server

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
